# py/app-market/market-analysis/crontab
10 0 28 * * cd /home/zhangyan.zyan/internal/market-analysis; ./run_pp.sh
